# weather-display.service
[Unit]
Description=E-Ink Weather Display
After=network.target

[Service]
Type=simple
User=austingibb
WorkingDirectory=/home/austingibb/tools/e-ink-pi-zero
ExecStart=/usr/bin/python3 /home/austingibb/tools/e-ink-pi-zero/code.py
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
