# weather-display.service
[Unit]
Description=E-Ink Weather Display
After=network.target

[Service]
Type=simple
User=austingibb
WorkingDirectory=/home/austingibb/tools/e-ink-pi-zero
EnvironmentFile=/home/austingibb/tools/e-ink-pi-zero/.env
ExecStart=/home/austingibb/tools/e-ink-pi-zero/venv/bin/python /home/austingibb/tools/e-ink-pi-zero/code.py
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
